# systemd/automation-prune.timer
[Unit]
Description=Daily Podman prune for automation stack

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
